# wiretide.service
[Unit]
Description=Wiretide Controller
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/wiretide
ExecStart=/opt/wiretide/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
